# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:cron]
command=/usr/sbin/cron -f
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
stdout_logfile=/var/log/cron.log
stderr_logfile=/var/log/cron-error.log

[program:gunicorn]
command=gunicorn groupplanning.wsgi:application --bind 0.0.0.0:8000 --workers 3 --timeout 60 --access-logfile - --error-logfile -
directory=/app
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
environment=PORT="8000",GUNICORN_WORKERS="3",GUNICORN_TIMEOUT="60"
stdout_logfile=/var/log/gunicorn.log
stderr_logfile=/var/log/gunicorn-error.log

[group:django]
programs=cron,gunicorn
